# systemd_units/elios-shadow-metrics.timer
[Unit]
Description=Run Elios Shadow Metrics daily at 03:57 +05 (22:57 UTC)

[Timer]
OnCalendar=*-*-* 22:57:00 UTC
AccuracySec=1min
Persistent=true
Unit=elios-shadow-metrics.service

[Install]
WantedBy=timers.target
