# systemd_units/elios-train-guard.service
[Unit]
Description=Elios Train Progress Guard
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/stockbot
Environment=PYTHONPATH=/root/stockbot
ExecStart=/root/stockbot/venv/bin/python -u /root/stockbot/tools/train_progress_guard.py
# WARN (exit=2) не считаем фейлом
SuccessExitStatus=0 2
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

[Install]
WantedBy=multi-user.target
