# systemd_units/elios-entry-gate.service
[Unit]
Description=Elios Entry Orders Gate
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
Environment=PYTHONUNBUFFERED=1
WorkingDirectory=/root/stockbot
EnvironmentFile=-/root/stockbot/.env
EnvironmentFile=-/root/stockbot/.env.local
ExecStart=/root/stockbot/venv/bin/python -m core.trading.guards.entry_orders_gate
Nice=5
IOSchedulingClass=2
