# systemd_units/elios-intraday-clamp.service
[Unit]
Description=Elios Intraday Clamp Guard
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
Environment=PYTHONUNBUFFERED=1
WorkingDirectory=/root/stockbot
EnvironmentFile=-/root/stockbot/.env
EnvironmentFile=-/root/stockbot/.env.local
ExecStart=/root/stockbot/venv/bin/python -m core.trading.guards.intraday_clamp
Nice=5
IOSchedulingClass=2

[Install]
WantedBy=multi-user.target
