# systemd_units/elios-data-quality.service
[Unit]
Description=Elios Data Quality Audit (dataset health)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/stockbot
Environment=PYTHONPATH=/root/stockbot
ExecStart=/root/stockbot/venv/bin/python -u /root/stockbot/core/pipeline/data_quality_report.py
# WARN (exit=2) не считаем ошибкой systemd
SuccessExitStatus=0 2
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

[Install]
WantedBy=multi-user.target
