# systemd_units/elios-loop.service
[Unit]
Description=Elios Trading Loop Launcher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/stockbot
Environment=PYTHONPATH=/root/stockbot
EnvironmentFile=-/root/stockbot/.env.local
ExecStart=/root/stockbot/venv/bin/python -u -m core.trading.loop_launcher
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Nice=5

[Install]
WantedBy=multi-user.target
