# systemd_units/eod-flatten.service
[Unit]
Description=Elios EOD Flatten (close all open positions near market close)
Wants=network-online.target
After=network-online.target
ConditionPathExists=/root/stockbot/core/trading/eod_flatten.py

[Service]
Type=oneshot
WorkingDirectory=/root/stockbot
Environment="PYTHONPATH=/root/stockbot"
EnvironmentFile=-/root/stockbot/.env.local
ExecStart=/root/stockbot/venv/bin/python3 -u -m core.trading.eod_flatten
User=root
Group=root
Nice=5

[Install]
WantedBy=multi-user.target
