# systemd_units/sell-health-probe.timer
[Unit]
Description=Run SELL health-probe periodically

[Timer]
OnBootSec=5min
OnUnitActiveSec=10min
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
