# systemd_units/elios-loop.timer.d/override.conf
[Timer]
OnCalendar=
OnCalendar=Mon..Fri *-*-* 13:30:00 UTC
Persistent=true
AccuracySec=1min
