# systemd_units/sell-notifier.service.d/00-env.conf
[Service]
WorkingDirectory=/root/stockbot
Environment=PYTHONPATH=/root/stockbot TZ=Asia/Tashkent
# Загружаем базовый и локальный env (локальный необязателен, но если есть — подхватится)
EnvironmentFile=/root/stockbot/.env
EnvironmentFile=-/root/stockbot/.env.local
