# systemd_units/elios-sell-once.service.d/01-cond-override.conf
[Service]
# Проверка окна торговли (скрипт должен вернуть 0 в окно)
ExecCondition=/root/stockbot/tools/within_sell_window.sh
