# systemd_units/elios-sell-once.service.d/98-prod.conf
[Service]
Environment=ELIOS_SELL_DRY_RUN=0
Environment=ELIOS_SELL_NO_POLL=0
Environment=ELIOS_SELL_MAX_WAIT_SEC=120
