# systemd_units/elios-sell-once.service.d/zz-override.conf
[Service]
# Стереть ВСЕ прежние условия (включая market_open_guard.py)
# Правильный гард: будни 18:30–23:59 и 00:00–00:30 со сдвигом на следующий день
