# systemd_units/eod-flatten.service.d/env.conf
[Service]
EnvironmentFile=/root/stockbot/.env.local
