# systemd_units/elios-loop.service.d/env.conf
[Service]
# load .env if present (optional with '-')
EnvironmentFile=-/root/stockbot/.env
