# systemd_units/elios-loop.service.d/mutex.conf
[Service]
# Переопределяем ExecStart, добавляя файловую блокировку
ExecStart=
ExecStart=/usr/bin/flock -n /run/elios.loop.mutex /root/stockbot/venv/bin/python -u -m core.trading.loop_launcher
