# systemd_units/elios-train-guard.service.d/override.conf
[Service]
Environment=ELIOS_SHADOW_MIN_P5=0.30
Environment=ELIOS_SHADOW_MIN_P10=0.25
Environment=ELIOS_SHADOW_MIN_P20=0.20
Environment=ELIOS_SHADOW_MAX_BRIER=0.22
